# misc/bigcake.socket
[Unit]
Description=bigcake socket

[Socket]
ListenStream=/run/bigcake.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target